# cron/crontab
* * * * * loxberry REPLACELBPBINDIR/watch.sh >> REPLACELBPLOGDIR/kidlox_cron.log 2>&1 
